# dist/zoondka-maps.service
[Unit]
Description=zoondka-maps app server
Documentation=https://github.com/zoondka/zoondka-maps
After=network.target

[Service]
Type=simple
PIDFile=%t/zoondka-maps.pid
User=zoondka-maps
Group=zoondka-maps
Environment="ZOONDKA-MAPS_PORT=16531"
ExecStart=/usr/bin/java -jar /srv/zoondka-maps/zoondka-maps.jar 16531
Restart=always
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=zoondka-maps

[Install]
WantedBy=multi-user.target
